# scripts/kiosk/smartphonecracy-kiosk.service
[Unit]
Description=Smartphonecracy fullscreen kiosk
Wants=network-online.target
After=network-online.target graphical-session.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
EnvironmentFile=/etc/smartphonecracy/kiosk.env
ExecStart=/opt/smartphonecracy/scripts/kiosk/launch-kiosk.sh
Restart=always
RestartSec=5
TimeoutStopSec=15
KillMode=mixed
SyslogIdentifier=smartphonecracy-kiosk

# Browser session runs as the dedicated, unprivileged account.
User=smartphonecracy-kiosk
Group=smartphonecracy-kiosk

[Install]
WantedBy=graphical.target
